# gitsync.service
[Unit]
Description=sync git

[Service]
ExecStart=/bin/bash /*path to sync.sh*/
